# deploy/wo-dashboard.service
# WordOps Dashboard - Systemd Service Unit
#
# Install to: /etc/systemd/system/wo-dashboard.service
#
# Commands:
#   systemctl daemon-reload
#   systemctl enable wo-dashboard
#   systemctl start wo-dashboard
#   systemctl status wo-dashboard

[Unit]
Description=WordOps Dashboard API
Documentation=https://github.com/wpbryant/WordOps-Dashboard
After=network.target
Wants=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root

# Working directory where the application is installed
WorkingDirectory=/var/www/wo-dashboard

# Load environment variables from config file
EnvironmentFile=/etc/wo-dashboard/config.env

# Start uvicorn on TCP port (simpler than unix socket for WordOps proxy)
ExecStart=/var/www/wo-dashboard/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000

# Create runtime directory
RuntimeDirectory=wo-dashboard
RuntimeDirectoryMode=0755

# Restart behavior
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:/var/log/wo-dashboard/access.log
StandardError=append:/var/log/wo-dashboard/error.log

# Security hardening disabled - WordOps needs full system access
# to manage apt, nginx configs, SSL certs, etc.
# NoNewPrivileges=false
# ProtectSystem=strict

[Install]
WantedBy=multi-user.target
